# docker/web/s6-rc.d/claude-watcher/finish
#!/command/with-contenv bash
# shellcheck shell=bash
# Exit 0 is normal here — watcher exits after each detected change and s6
# restarts it immediately. Only log unexpected non-zero exits.
exit_code="${1:-?}"
signal="${2:--1}"

if [ "$exit_code" -ne 0 ] && [ "$signal" -eq -1 ]; then
    echo "[claude-watcher] UNEXPECTED EXIT — code=${exit_code} — s6 will restart" >&2
fi
