# docker/web/s6-rc.d/claude-watcher/run
#!/command/with-contenv bash
# shellcheck shell=bash
# Watches config paths and restarts claude-session so new agents, skills, hooks,
# and config edits are picked up without a full container restart.
# Uses inotifywait (inotify-tools) — one-shot per change, then exits.
# s6 restarts this watcher immediately, making the watch continuous.
#
# Watched events:
#   modify   — file content changed (edit in place)
#   create   — new file/dir added (new agent, skill, or command)
#   delete   — file/dir removed
#   moved_to — file renamed/moved into a watched dir
#   move_self — watched file itself was moved

workspace="${AXON_WORKSPACE:-/workspace}"
home_dir="${HOME:-/home/node}"

# Whitelist of paths that affect Claude's runtime behavior.
# Intentionally narrow — avoids watching runtime-written dirs like
# ~/.claude/projects/ (session history) or ~/.claude/statsig/ (analytics)
# which would cause a restart loop while Claude is running.
#
# Dirs are watched recursively (-r); files are watched directly (-f).
# Each entry is only added if the path currently exists.

add() {
    local flag="$1" path="$2"
    if [ "$flag" = "-r" ] && [ -d "$path" ]; then
        watch_args+=(-r "$path")
    elif [ "$flag" = "-f" ] && [ -f "$path" ]; then
        watch_args+=("$path")
    fi
}

watch_args=()

# ── User-level (~/.claude/) ────────────────────────────────────────────────
add -r "${home_dir}/.claude/agents"
add -r "${home_dir}/.claude/commands"
add -r "${home_dir}/.claude/hooks"
add -r "${home_dir}/.claude/output-styles"
add -r "${home_dir}/.claude/plugins"
add -r "${home_dir}/.claude/skills"
add -f "${home_dir}/.claude/settings.json"
add -f "${home_dir}/.claude/settings.local.json"
# ~/.claude.json intentionally excluded — Claude writes auth/session state here
# at runtime, which would create a restart loop if watched.

# ── Project-level (${workspace}/.claude/) ─────────────────────────────────
add -r "${workspace}/.claude/agents"
add -r "${workspace}/.claude/commands"
add -r "${workspace}/.claude/hooks"
add -r "${workspace}/.claude/output-styles"
add -f "${workspace}/.claude/settings.json"
add -f "${workspace}/.claude/settings.local.json"
add -f "${workspace}/CLAUDE.md"
add -f "${workspace}/.mcp.json"

if [ ${#watch_args[@]} -eq 0 ]; then
    echo "[claude-watcher] no watchable paths found under ${home_dir} or ${workspace} — retrying in 30s" >&2
    sleep 30
    exit 0
fi

echo "[claude-watcher] watching: ${watch_args[*]}" >&2

# Block until any change is detected, then exit (s6 restarts us immediately).
changed="$(inotifywait -e modify,create,delete,moved_to,move_self \
    "${watch_args[@]}" \
    --format '%w%f' \
    --quiet 2>/dev/null)"

echo "[claude-watcher] change detected: ${changed} — restarting claude-session" >&2

# 1-second debounce: absorbs rapid multi-file changes (e.g. saving an agent
# file triggers modify + moved_to). Without this, the session would restart
# mid-write and potentially load a half-written file.
sleep 1

s6-svc -r /run/service/claude-session
