# docker/web/s6-rc.d/claude-session/finish
#!/command/with-contenv bash
# shellcheck shell=bash
# $1 = exit code, $2 = signal (-1 if clean exit, 15 = SIGTERM from watcher restart)
exit_code="${1:-?}"
signal="${2:--1}"

if [ "$signal" -eq 15 ]; then
    echo "[claude-session] restarted by watcher (SIGTERM) — reloading config" >&2
elif [ "$exit_code" -ne 0 ] || [ "$signal" -ne -1 ]; then
    echo "[claude-session] UNEXPECTED EXIT — code=${exit_code} signal=${signal} — s6 will restart" >&2
else
    echo "[claude-session] clean exit — code=${exit_code}" >&2
fi
