# docker/web/s6-rc.d/claude-session/run
#!/command/with-contenv bash
# shellcheck shell=bash
# Persistent Claude Code session running in the web container.
# Runs in the mounted workspace so CLAUDE.md and project .claude/ are loaded.
# --continue resumes the most recent conversation; --fork-session avoids session
# ID collisions on rapid restarts triggered by claude-watcher.
#
# --dangerously-skip-permissions: required because s6 services run without a TTY.
# Claude detects no TTY and refuses to start in interactive mode unless this flag
# is set. This is intentional — the container is an isolated sandbox.
#
# script -q -e /dev/null: allocates a pseudo-TTY so Claude runs in interactive
# (REPL) mode rather than exiting immediately. -e propagates Claude's exit code.
# /dev/null discards the script session transcript.
cd "${AXON_WORKSPACE:-/workspace}"
exec s6-setuidgid node \
    script -q -e /dev/null \
        -c 'claude --continue --fork-session --dangerously-skip-permissions'
